# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:myclock]
command=python3 /app/my_clock.py
autostart=true
autorestart=true